# k8/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://proxy.sgp.hp.com:8080" "NO_PROXY=localhost,127.0.0.1,https://registry-1.docker.io"
